# supervisord.conf
[supervisord]
# The path to the activity log of the supervisord process. 
# This option can include the value %(here)s, which expands to the 
# directory in which the supervisord configuration file was found.
logfile=/home/budi/supervisord.log

# The maximum number of bytes that may be consumed by the activity log file 
# before it is rotated (suffix multipliers like “KB”, “MB”, and “GB” can be 
# used in the value). 
# Set this value to 0 to indicate an unlimited log size.
logfile_maxbytes=50MB

# The number of backups to keep around resulting from 
# activity log file rotation. 
# If set to 0, no backups will be kept.
logfile_backups=10

# The logging level, dictating what is written to the supervisord activity log.
# One of critical, error, warn, info, debug, trace, or blather.
# Note that at log level debug, the supervisord log file will record 
# the stderr/stdout output of its child processes and extended info info 
# about process state changes, which is useful for debugging a process 
# which isn’t starting properly. 
loglevel=info

# The location in which supervisord keeps its pid file.
# This option can include the value %(here)s, which expands 
# to the directory in which the supervisord configuration file was found.
pidfile=/home/budi/supervisord.pid

# The umask of the supervisord process.
umask=022

# If true, supervisord will start in the foreground instead of daemonizing.
nodeamon=false

# The minimum number of file descriptors that must be available before 
# supervisord will start successfully. 
# A call to setrlimit will be made to attempt to raise the soft and hard limits
# of the supervisord process to satisfy minfds.
# The hard limit may only be raised if supervisord is run as root. 
# supervisord uses file descriptors liberally, and will enter a failure mode 
# when one cannot be obtained from the OS, so it’s useful to be able to 
# specify a minimum value to ensure it doesn’t run out of them during 
# execution. 
# This option is particularly useful on Solaris, which has a low per-process 
# fd limit by default.
minfds=1024

# The minimum number of process descriptors that must be available 
# before supervisord will start successfully. A call to setrlimit will be 
# made to attempt to raise the soft and hard limits of the supervisord 
# process to satisfy minprocs. The hard limit may only be raised if 
# supervisord is run as root. supervisord will enter a failure mode 
# when the OS runs out of process descriptors, so it’s useful to 
# ensure that enough process descriptors are available upon 
# supervisord startup.
minprocs=200

# Prevent supervisord from clearing any existing AUTO chlild log files 
# at startup time. Useful for debugging.
nocleanup=false

# The directory used for AUTO child log files. This option can include 
# the value %(here)s, which expands to the directory in which the 
# supervisord configuration file was found.
# childlogdir=

# If supervisord is run as the root user, switch users to this UNIX user 
# account before doing any meaningful processing. This value has no effect 
# if supervisord is not run as root. Default value = current user
#user=

# When supervisord daemonizes, switch to this directory. This option 
# can include the value %(here)s, which expands to the directory in which 
# the supervisord configuration file was found. 
# Default value = current directory
#directory=

# Strip all ANSI escape sequences from child log files.
strip_ansi=false

# A list of key/value pairs in the form KEY="val",KEY2="val2" that will 
# be placed in the supervisord process’ environment (and as a result in 
# all of its child process’ environments). This option can include the 
# value %(here)s, which expands to the directory in which the supervisord 
# configuration file was found. Values containing non-alphanumeric 
# characters should be quoted (e.g. KEY="val:123",KEY2="val,456"). 
# Otherwise, quoting the values is optional but recommended. 
# To escape percent characters, simply use two. (e.g. URI="/first%%20name") 
# Note that subprocesses will inherit the environment variables of the 
# shell used to start supervisord except for the ones overridden here 
# and within the program’s environment option. Default value = no value
#environment=

# The identifier string for this supervisor process, used by the RPC 
# interface.
identifier=supervisor

[supervisorctl]
# The URL that should be used to access the supervisord server, 
# e.g. http://localhost:9001. For UNIX domain sockets, use 
# unix:///absolute/path/to/file.sock.
serverurl=http://127.0.0.1:9001

# The username to pass to the supervisord server for use in authentication. 
# This should be same as username from the supervisord server configuration 
# for the port or UNIX domain socket you’re attempting to access.
username=budi

# The password to pass to the supervisord server for use in authentication. 
# This should be the cleartext version of password from the supervisord 
# server configuration for the port or UNIX domain socket you’re attempting 
# to access. This value cannot be passed as a SHA hash. Unlike other 
# passwords specified in this file, it must be provided in cleartext.
password=budibodonk

# String used as supervisorctl prompt.
prompt=supervisor

# A path to use as the readline persistent history file. If you enable 
# this feature by choosing a path, your supervisorctl commands will be 
# kept in the file, and you can use readline (e.g. arrow-up) to invoke 
# commands you performed in your last supervisorctl session.
#history_file=

[inet_http_server]
port=127.0.0.1:9001
username=budi
password=budibodonk

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:foo]
command=/bin/cat
#process_name=%(program_name)s
#numprocs=1
#directory=/tmp
#umask=022
#priority=999
#autostart=true
#autorestart=true
#startsecs=10
#startretries=3
#exitcodes=0,2
#stopsignal=TERM
#stopwaitsecs=10
#user=chrism
#redirect_stderr=false
#stdout_logfile=/a/path
#stdout_logfile_maxbytes=1MB
#stdout_logfile_backups=10
#stdout_capture_maxbytes=1MB
#stderr_logfile=/a/path
#stderr_logfile_maxbytes=1MB
#stderr_logfile_backups=10
#stderr_capture_maxbytes=1MB
#environment=A="1",B="2"
#serverurl=AUTO
